# Web/BlindAsABat/supervisord.conf
[supervisord]
nodaemon=true

[program:blindsqli]
command=gunicorn -b 0.0.0.0:5000 -w 8 app:app 0.0.0.0 5000

[program:mysqld]
command=mysqld
